# deploy/systemd/cco-backend.service
[Unit]
Description=CCO backend (FastAPI + uvicorn)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/cco/backend
ExecStart=/srv/cco/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
